# ws_light.service
[Unit]
Description=WS Light

[Service]
Type=simple
ExecStart=/usr/local/bin/ws_light

[Install]
WantedBy=multi-user.target
